# templates/rc.local.j2


{% if GraphicusMaximus.enable %}
cd /opt/services/GraphicusMaximus && sudo -u graphicus screen -dmS GraphicusMaximus java -jar graphicus_maximus.jar 9988 'jdbc:mysql://localhost/graphicus?user=root&password=blar'
{% endif %}

{% if TicTacToe.enable %}
sudo -u tictactoe screen -dmS TicTacToe_userd /opt/services/TicTacToe/run.sh userd.rb
sudo -u tictactoe screen -dmS TicTacToe_userd /opt/services/TicTacToe/run.sh tttd.rb 
{% endif %}

{% if Wizard.enable %}
sudo -u wizard screen -dmS Wizard java -jar /opt/services/Wizard/wizard.jar
{% endif %}

{% if PythonOnPills.enable %}
cd {{ctfc_service_dir}}/PythonOnPills; python pills.py &
{% endif %}


exit 0
